# sci-ml/lemonade/files/lemond.initd
#!/sbin/openrc-run
# Copyright 1999-2026 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="Lemonade local LLM server"

command="/usr/bin/lemond"
command_args="${LEMOND_OPTS}"
command_user="lemonade:lemonade"
directory="/var/lib/lemonade"
supervisor="supervise-daemon"
output_log="/var/log/lemonade/lemond.log"
error_log="/var/log/lemonade/lemond.log"

depend() {
	need net
	use dns
}

start_pre() {
	checkpath -d -o lemonade:lemonade -m 0750 /var/lib/lemonade
	checkpath -d -o lemonade:lemonade -m 0750 /var/log/lemonade
	checkpath -d -o lemonade:lemonade -m 0700 /run/lemonade
}
