# sci-biology/foldingathome/files/foldingathome-initd
#!/sbin/openrc-run
# Copyright 1999-2026 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="Folding@home Client"

command="/opt/foldingathome/fah-client"
command_args="--config=/etc/fah-client/config.xml --log=/var/log/fah-client/log.txt --log-rotate-dir=/var/log/fah-client/ ${FAH_OPTS}"
command_user="foldingathome:foldingathome"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

directory="/var/lib/fah-client"

depend() {
	need net
	after logger
}
