# net-proxy/snowflake/files/snowflake-proxy.initd
#!/sbin/openrc-run
# Copyright 2023-2026 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="WebRTC pluggable transport proxy for Tor"

command="/usr/bin/snowflake-proxy"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need net
	after dns
}
